# config/.config/systemd/user/plasma-krunner.service.d/oom.conf
# Prevent an OOM kill of a child from killing every process in this unit (eg. Emacs server)
[Service]
OOMPolicy=continue
